# systemd/user/copilot-api.service
[Unit]
Description=copilot-api (Anthropic-native proxy for GitHub Copilot, for Claude Code)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/usr/bin:/bin
ExecStart=/home/linuxbrew/.linuxbrew/bin/copilot-api start --port 4001 --account-type enterprise
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
