# home/.config/systemd/user/caddy.service
[Unit]
Description=Caddy HTTPS reverse proxy (litellm + copilot-api)
After=network.target

[Service]
ExecStart=/home/linuxbrew/.linuxbrew/bin/caddy run --config /home/xingfeixu/.config/caddy/Caddyfile
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
