# home/.config/systemd/user/code-serve-web.service
[Unit]
Description=VS Code Server (official serve-web, bound to localhost for tailscale serve)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/xingfeixu/.local/bin/code serve-web --host 127.0.0.1 --port 4321 --without-connection-token --accept-server-license-terms
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
